# admin/setup/backup-remote.service
[Unit]
Description=Download backup of thefiletree.com

[Service]
Type=oneshot
ExecStart=/home/tyl/file/tree/admin/setup/backup-remote
