# admin/setup/backup-remote.timer
[Unit]
Description=Download backup of thefiletree.com

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
