# admin/setup/renew-cert.service
[Unit]
Description=Update thefiletree.com TLS certificate

[Service]
Type=oneshot
ExecStart=/home/dom/tree/admin/setup/renew-cert
