# admin/setup/renew-cert.timer
[Unit]
Description=Run thefiletree.com certificate renewal daily

[Timer]
OnCalendar=6/12:42
