# admin/setup/tree.service
[Unit]
Description=TheFileTree server

[Service]
ExecStart=/home/dom/tree/admin/setup/tree
Restart=on-failure

[Install]
WantedBy=multi-user.target
